# .config/systemd/user/auto-ssh-config.service
[Install]
WantedBy=default.target

[Service]
ExecStart=%h/.local/lib/dotfiles/auto-ssh-config

[Unit]
Description=auto-ssh-config
